# ble2mqtt.service
[Unit]
Description=A naive mimic of zigbee2mqtt for bluetooth with python
After=network.target bluetooth.target

[Service]
DynamicUser=yes
BindReadOnlyPaths=/opt/ble2mqtt:/ble2mqtt
BindReadOnlyPaths=/opt/ble2mqtt/config:/ble2mqtt/config
WorkingDirectory=/ble2mqtt
ExecStart=/bin/python3 /ble2mqtt/main.py
Restart=on-failure
RestartSec=20

[Install]
WantedBy=multi-user.target
